# neptun-daemon.service
[Unit]
Description=NEPTUN Telegram Data Fetcher Daemon
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/neptun
Environment=PATH=/var/www/neptun/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/var/www/neptun/venv/bin/python3 /var/www/neptun/neptun_daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/neptun_daemon.log
StandardError=append:/var/log/neptun_daemon.log

[Install]
WantedBy=multi-user.target
